# pi/sensor-reader/sensor-reader-supervisor.conf
[program:weerstation-sensor-reader]
command=python /home/pi/weerstation.pi/sensor-reader/serial_data_reader.py
autostart=true
autorestart=true
user=root
